# htb-neurogrid-2025/secure/sakuras_embrace/config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770
chown=root:editor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[supervisord]
user=root
nodaemon=true
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:node]
command=node app.js
user=challenger
directory=/app/main_app
autostart=true
autorestart=true
stdout_logfile=/var/log/node/access.log
stderr_logfile=/var/log/node/error.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true


[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile --resume
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/caddy/access.log
stderr_logfile=/var/log/caddy/error.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

;[program:editor-app]
;command=python3 app.py
;user=editor
;directory=/app/editor/backend
;autostart=true
;autorestart=true
;stdout_logfile=/app/editor/backend/1.out
;stderr_logfile=/app/editor/backend/1.err
;stdout_logfile_maxbytes=0
;stderr_logfile_maxbytes=0